# myapp.ini.j2
[program:myapp]
priority=10
directory=/tmp/myapp
command=python app.py
user=root
autostart=true
autorestart=true
